# muster.service
[Unit]
Description=Muster MCP Server
After=network.target
Requires=muster.socket
After=muster.socket

[Service]
Type=simple
ExecStart="%h/go/bin/muster" serve --debug --yolo --config-path "%h/.config/muster"
Restart=always
RestartSec=5
Environment=PATH=%h/go/bin:/usr/local/bin:/usr/bin:/bin

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=muster

# Security
NoNewPrivileges=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictRealtime=true

[Install]
WantedBy=default.target 